# templates/default_systemctl.timer
[Unit]
Description=%%app_title%% Backup Scheduling.

[Timer]
OnCalendar=%%frequency%%
Persistent=true

[Install]
WantedBy=multi-user.target
